# docker/supervisord/php-fpm.conf
[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_maxbytes=0
stderr_logfile=/dev/stdout
stderr_maxbytes=0
exitcodes=0